# supervisord.conf
; EuraFlow Supervisor 配置文件
; 项目级配置，不需要系统权限

[unix_http_server]
file=%(here)s/tmp/supervisor.sock   ; socket 文件位置
chmod=0700                           ; socket 文件权限

[supervisord]
logfile=%(here)s/logs/supervisord.log ; 主日志文件
logfile_maxbytes=50MB                 ; 日志文件最大大小
logfile_backups=10                    ; 日志备份数量
loglevel=info                         ; 日志级别
pidfile=%(here)s/tmp/supervisord.pid  ; PID 文件
nodaemon=false                        ; 后台运行
minfds=1024                          ; 最小文件描述符
minprocs=200                         ; 最小进程数
directory=%(here)s                    ; 工作目录

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/tmp/supervisor.sock ; 使用 UNIX socket

; 项目进程组
[group:euraflow]
programs=backend,worker,celery_worker,celery_beat
priority=999

; 后端 API 服务
[program:backend]
command=%(here)s/venv/bin/uvicorn ef_core.app:app --host 0.0.0.0 --port 8000 --log-level info
directory=%(here)s
autostart=true
autorestart=true
startsecs=2
stopwaitsecs=10
stdout_logfile=%(here)s/logs/backend.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/logs/backend-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PYTHONPATH="%(here)s",PATH="%(here)s/venv/bin:%(ENV_PATH)s"
priority=100

; 水印任务处理器
[program:worker]
command=%(here)s/venv/bin/python -m plugins.ef.channels.ozon.services.watermark_task_runner
directory=%(here)s
autostart=true
autorestart=true
startsecs=2
stopwaitsecs=10
stdout_logfile=%(here)s/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/logs/worker-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PYTHONPATH="%(here)s",PATH="%(here)s/venv/bin:%(ENV_PATH)s"
priority=200

; 通用 Celery Worker（处理异步任务，如标签下载等）
[program:celery_worker]
command=%(here)s/venv/bin/celery -A ef_core.tasks.celery_app worker --loglevel=info --queues=celery,default,ef_pull,ef_push,ef_core --pool=prefork --concurrency=2 --max-memory-per-child=150000
directory=%(here)s
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=10
stdout_logfile=%(here)s/logs/celery-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/logs/celery-worker-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PYTHONPATH="%(here)s",PATH="%(here)s/venv/bin:%(ENV_PATH)s"
priority=150

; Celery Beat（定时任务调度器）
[program:celery_beat]
command=%(here)s/venv/bin/celery -A ef_core.tasks.celery_app beat --loglevel=info --scheduler=celery.beat:PersistentScheduler
directory=%(here)s
autostart=true
autorestart=true
startsecs=2
stopwaitsecs=10
stdout_logfile=%(here)s/logs/celery-beat.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/logs/celery-beat-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PYTHONPATH="%(here)s",PATH="%(here)s/venv/bin:%(ENV_PATH)s"
priority=160

; Web 管理界面（已禁用，使用命令行工具 supervisorctl）
; 如需启用，请使用环境变量配置密码：
; [inet_http_server]
; port=127.0.0.1:9001
; username=admin
; password=%(ENV_SUPERVISOR_PASSWORD)s