# deploy/systemd/euraflow-backend.service
[Unit]
Description=EuraFlow Backend API Service
Documentation=https://github.com/your-org/EuraFlow
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service
Wants=network-online.target

[Service]
Type=notify
User=euraflow
Group=euraflow
WorkingDirectory=/opt/euraflow

# 环境变量
EnvironmentFile=/opt/euraflow/.env
Environment="PATH=/opt/euraflow/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONPATH=/opt/euraflow"

# 启动前检查
ExecStartPre=/opt/euraflow/venv/bin/python -c "import sys; sys.exit(0)"
ExecStartPre=/opt/euraflow/venv/bin/alembic upgrade head

# 主进程
ExecStart=/opt/euraflow/venv/bin/gunicorn ef_core.app:app \
    --bind 0.0.0.0:8000 \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --worker-connections 1000 \
    --max-requests 10000 \
    --max-requests-jitter 1000 \
    --timeout 60 \
    --keepalive 5 \
    --graceful-timeout 30 \
    --access-logfile /var/log/euraflow/access.log \
    --error-logfile /var/log/euraflow/error.log \
    --log-level info \
    --capture-output \
    --enable-stdio-inheritance \
    --statsd-host=localhost:8125 \
    --statsd-prefix=euraflow.backend

# 重载配置
ExecReload=/bin/kill -s USR2 $MAINPID

# 停止服务
ExecStop=/bin/kill -s TERM $MAINPID
KillMode=mixed
KillSignal=SIGQUIT
TimeoutStopSec=30

# 重启策略
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/euraflow /var/log/euraflow /var/run/euraflow

# 性能调优
CPUAccounting=true
MemoryAccounting=true
MemoryMax=2G
TasksMax=100

# 日志
StandardOutput=append:/var/log/euraflow/backend.log
StandardError=append:/var/log/euraflow/backend-error.log

[Install]
WantedBy=multi-user.target